# systemd-services-backup/household_agent.service
[Unit]
Description=Household Agent - домашній асистент
After=network.target

[Service]
Type=simple
User=sashok
WorkingDirectory=/home/sashok/.openclaw/workspace/household_agent
ExecStartPre=/bin/mkdir -p /run/user/1000/bot-locks
ExecStart=/usr/bin/flock -n /run/user/1000/bot-locks/household_agent.lock /home/sashok/.openclaw/workspace/household_agent/venv/bin/python3 main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
