# systemd-services-backup/sam.service
[Unit]
Description=Sam — Personal Agent
After=network.target

[Service]
Type=simple
User=sashok
WorkingDirectory=/home/sashok/.openclaw/workspace/sam
EnvironmentFile=/home/sashok/.openclaw/workspace/sam/.env
ExecStartPre=/bin/mkdir -p /run/user/1000/bot-locks
ExecStart=/usr/bin/flock -n /run/user/1000/bot-locks/sam.lock /home/sashok/.openclaw/workspace/sam/venv/bin/python main.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
